# rocks_hpc_init/post_boot_script.sh
#!/bin/bash

#chkconfig: 2345 99 30
#description: starts post commands
#processname: post_boot_script

